# systemd/claude-auto-retry-reconcile.timer
[Unit]
Description=Periodically reconcile claude-auto-retry monitors (self-healing coverage)
Documentation=https://github.com/cheapestinference/claude-auto-retry

[Timer]
# First run shortly after the user session comes up (lets tmux settle)...
OnStartupSec=2min
# ...then every 5 min. A monitor that dies is re-armed within one interval; a
# crash mid-window loses at most ~5 min of coverage on that pane.
OnUnitActiveSec=5min
# NOTE: no Persistent= here — it only affects OnCalendar= (wall-clock) timers. These are
# monotonic (OnStartupSec/OnUnitActiveSec), which pause across suspend and resume on wake,
# so there is nothing for Persistent to "catch up".

[Install]
WantedBy=timers.target
